# test_files/etc/supervisor/conf.d/run_for_ten.conf
[program:runforten]
command=bash test_files/bin/run_for_ten.sh
autostart=true
autorestart=true
startretries=3
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
